# provisioning/roles/mailhog/templates/systemd/mailhog.service.j2
[Unit]
Description=
Documentation=
User={{ mailhog_user }}
Group={{ mailhog_group }}

[Service]
ExecStart={{ mailhog_home_dir }}/bin/mailhog

[Install]
WantedBy=multi-user.target
